# templates/gateway_dl_config.conf.erb
description "Copy ec2-user-data to /etc/default/"

start on started ec2_user_data and filesystem

console log

script
    echo "Preparing gateway config"
    /usr/local/bin/gateway_dl_config
    echo "Restarting nginx"
    service nginx restart
    echo "Completed"
end script