# templates/gateway.conf.erb
description "Gateway Service"

start on started gateway_dl_config and stopped ec2_user_data

console log

limit nofile <%= @gateway_limit_nofile %> <%= @gateway_limit_nofile %>

script

  echo "Gateway: performaing startup steps"

  #Make sure all required info are present in ec2-user-data.
  echo "Making sure ec2-user-data contains required definitions (Wait until it becomes available)"
  EC2File='/etc/default/ec2-user-data'
  while ! grep -q EC2_REGION "$EC2File"; do
    echo -n "."
    sleep 1
  done
  echo " :Done"

  . /etc/default/ec2-user-data

  URL="http://169.254.169.254/latest/"
  ID=$(curl -s $URL/meta-data/instance-id)
  STACK_ID=$(aws ec2 describe-tags --region ${EC2_REGION} --filter Name=resource-id,Values=${ID} | jq -r .Tags | jq '.[] | select(.Key=="aws:cloudformation:stack-id")' | jq ".Value")

  echo "Instance ID : ${ID}"
  echo "Stack ID : ${STACK_ID}"

  while ! curl -s -k -o /dev/null -I -w "%{http_code}" https://localhost/sys/health | grep -q '200' ;
  do
   echo 'Waiting for gateway services to come alive'
   sleep 1
  done
  echo 'Gateway services are healthy'

  # Single to AWS that we are up
  echo 'Sending Gateway CF signal'
  aws --region ${EC2_REGION} cloudformation signal-resource --logical-resource-id GatewayAutoScalingGroup --stack-name $(echo "$STACK_ID" | sed -e 's/^"//' -e 's/"$//') --status SUCCESS --unique-id ${ID}
  echo 'Completed'

end script